# systemd/service/fixity_comparison.service
[Unit]
Description=Archival Storage Fixity Comparison Server

[Service]
Type=simple
WorkingDirectory=/cul/app/archival_storage_ingest
ExecStart=/cul/app/archival_storage_ingest/scripts/fixity_comparison.sh
Restart=on-failure
PrivateTmp=true

User=cular
Group=cular

[Install]
WantedBy=multi-user.target